# com.user.vpn-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.vpn-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__SCRIPT_PATH__</string>
    </array>

    <!-- 每 N 秒執行一次（由 CHECK_INTERVAL 設定） -->
    <key>StartInterval</key>
    <integer>__CHECK_INTERVAL__</integer>

    <!-- 開機後立即執行一次 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- 持續運行（不要因為一次失敗就停止） -->
    <key>KeepAlive</key>
    <false/>

    <!-- 標準輸出/錯誤重導向到日誌 -->
    <key>StandardOutPath</key>
    <string>/tmp/vpn_monitor_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/vpn_monitor_stderr.log</string>

    <!-- 環境變數 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
